# systemd/.config/systemd/user/mailsync.timer
[Unit]
Description=mailsync timer

[Timer]
OnCalendar=*:0/3
Unit=mailsync.service

[Install]
WantedBy=default.target